# docker/queue-worker.conf
[program:laravel-worker]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/worker.log
priority=3